# veripay-bot.service
[Unit]
Description=VeriPay Telegram Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=macbook
Group=staff
WorkingDirectory=/Users/macbook/veripay
ExecStart=/usr/bin/python3 /Users/macbook/veripay/lean_veripay_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=veripay-bot

# Environment variables
Environment=PYTHONPATH=/Users/macbook/veripay
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/Users/macbook/veripay/logs /Users/macbook/veripay/uploads /Users/macbook/veripay/reports

[Install]
WantedBy=multi-user.target 